# __APP_SLUG__/rootfs/etc/s6-overlay/s6-rc.d/init-__APP_SLUG__/up
/etc/s6-overlay/s6-rc.d/init-__APP_SLUG__/run
